# cron.conf
[program:cron]
command=sudo cron -f -L 15
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cron.stderr.log
stdout_logfile=/var/log/supervisor/cron.stdout.log
